# server/open-reading-web/deploy/open-reading-web.service
[Unit]
Description=Open Reading website and release service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=open-reading
Group=open-reading
WorkingDirectory=/srv/open-reading/current
EnvironmentFile=/srv/open-reading/shared/.env
ExecStart=/srv/open-reading/current/.venv/bin/uvicorn app.main:create_app --factory --host 127.0.0.1 --port 3002 --proxy-headers --forwarded-allow-ips=127.0.0.1 --no-access-log
Restart=on-failure
RestartSec=3
TimeoutStartSec=30
TimeoutStopSec=20
KillSignal=SIGINT
LimitNOFILE=65536
MemoryMax=256M

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
ReadWritePaths=/srv/open-reading/data /srv/open-reading/releases

[Install]
WantedBy=multi-user.target
