# server/open-reading-web/deploy/open-reading-backup.timer
[Unit]
Description=Daily Open Reading release metadata backup

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true
RandomizedDelaySec=20m

[Install]
WantedBy=timers.target

