# server/open-reading-web/deploy/open-reading-backup.service
[Unit]
Description=Back up Open Reading release metadata

[Service]
Type=oneshot
User=open-reading
Group=open-reading
EnvironmentFile=/srv/open-reading/shared/.env
ExecStart=/srv/open-reading/current/deploy/backup.sh
Nice=10
IOSchedulingClass=idle
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/open-reading/backups
ReadOnlyPaths=/srv/open-reading/data

